# jetbot.service
[Unit]
Description=Jetbot Robot Launcher Service
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=jetbot
Group=jetbot
WorkingDirectory=/opt/jetbot
ExecStart=/usr/bin/python3 /opt/jetbot/jetbot_launcher.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/jetbot/service.log
StandardError=append:/var/log/jetbot/service_error.log

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="JETBOT_HOME=/opt/jetbot"

# Resource limits (optional - adjust based on your needs)
MemoryLimit=2G
CPUQuota=80%

# Restart policy
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
